# raspberryfluke.service
[Unit]
Description=RaspberryFluke network discovery display service
# Use network-pre.target instead of network.target so the service starts
# as soon as the network interface exists, without waiting for a full
# network stack (DHCP lease, DNS, etc.) to be ready. RaspberryFluke
# manages its own DHCP awareness via /proc/net/route polling.
After=sysinit.target local-fs.target
Wants=network-pre.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/raspberryfluke
Environment=PYTHONPATH=/opt/raspberryfluke
ExecStart=/usr/bin/python3 /opt/raspberryfluke/main.py
Restart=on-failure
RestartSec=3
TimeoutStopSec=15

# Route all output to the systemd journal.
StandardOutput=journal
StandardError=journal

# Systemd hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
LockPersonality=true
RestrictSUIDSGID=true
RemoveIPC=true

# AF_PACKET  is required for raw socket capture (capture_raw.py, trigger.py).
# AF_UNIX    is required for general Python runtime communication.
# AF_NETLINK is required for network interface queries.
# AF_INET    is required for SNMP UDP queries and DHCP socket operations.
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_PACKET AF_INET

[Install]
WantedBy=multi-user.target
