# raspberryfluke.service
[Unit]
Description=RaspberryFluke network discovery display service
After=network.target lldpd.service
Requires=lldpd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/raspberryfluke
Environment=PYTHONPATH=/opt/raspberryfluke
ExecStart=/usr/bin/python3 /opt/raspberryfluke/main.py
Restart=on-failure
RestartSec=3

# Systemd hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
LockPersonality=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target